# etc/gotham/etc/init/gotham.conf
# gotham - warpsmith comics HTTP API

description "warpsmith comics HTTP API"

start on stopped rc
respawn

exec /usr/local/gotham/gotham
